# centos/openresty/openresty.conf
[program:openresty]
command=/usr/local/openresty/nginx/sbin/nginx -p /app -c conf/nginx.conf  -g 'daemon off;'
; autorestart=true