# centos/openresty/openresty.conf
[program:openresty]
command=/usr/local/openresty/nginx/sbin/nginx -p /var/www/ -c config/nginx.conf  -g 'daemon off;'
