# centos/tomcat7/tomcat.conf
[program:tomcat]
directory=/tomcat
command=/tomcat/bin/catalina.sh run