# hive/supervisor_conf/supervisor_hive.conf
[program:metastore]
command=/etc/supervisor/conf.d/start_hive_metastore.sh

[program:hiveserver2]
command=/etc/supervisor/conf.d/start_hive_hiveserver2.sh
